# alerter.cron
* * * * * /Users/DPurkayastha/.nvm/versions/node/v16.14.0/bin/node /Users/DPurkayastha/code/misc/mine/utils/my-mtg-alerter/mtgScript.js nowish-notify
